# proxy-injector.service
# proxy-injector.service
# Systemd service file for proxy-injector
# Copy this to /etc/systemd/system/proxy-injector.service

[Unit]
Description=Proxy Injector - Reverse proxy with nonce script injection
Documentation=https://github.com/minhqnd/proxy-injector
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/proxy-injector
ExecStart=/usr/bin/node proxy-inject-nonce.js
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=proxy-injector

# Environment variables
Environment=NODE_ENV=production
Environment=PORT=3000

# Security settings
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/home/ubuntu/proxy-injector
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes

# Resource limits
LimitNOFILE=65535
LimitNPROC=4096

[Install]
WantedBy=multi-user.target